# ddns-dhcp-sync/files/ddns-dhcp-sync.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/ddns-dhcp-sync.sh --noexit
	procd_close_instance
}

service_triggers()
{
        procd_add_reload_trigger "dhcp" "ddns-dhcp-sync"
}
